# resources/linux/servant.service
[Unit]
Description=Foo

[Service]
ExecStart=/usr/bin/servant runserver
User=servant
Group=servant

[Install]
WantedBy=multi-user.target